# packaging/pym2.service
[Unit]
Description=Pym2 Python Process Manager Agent
After=network.target

[Service]
Type=simple
User=pym2
Group=pym2
ExecStart=/usr/bin/pym2 agent
Restart=always
RestartSec=2
RuntimeDirectory=pym2
RuntimeDirectoryMode=0755
StateDirectory=pym2
StateDirectoryMode=0755
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
